# mdstack-service_2.2/mdstack-service_2/mdstack-service_2.0/service/datastats.service
[Unit]
Description=Data Statistics
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/datastats_service.py

[Install]
WantedBy=multi-user.target
