# mdstack-service_2.2/mdstack-service_2/mdstack-service_2.0/service/pullfiles.service
[Unit]
Description=Pull Files
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/pullfiles_service.py

[Install]
WantedBy=multi-user.target
