# mdstack-service_2.2/mdstack-service_2/mdstack-service_2.0/service/taskscheduler.service
[Unit]
Description=Task Scheduler
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/taskscheduler_service.py

[Install]
WantedBy=multi-user.target
